# deploy/express-bot.service
[Unit]
Description=Express Samar Bot
After=network-online.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Папка з кодом (deзмініть на реальний шлях)
WorkingDirectory=/home/ubuntu/express_x016

# Віртуальне середовище (deзмініть на реальний шлях)
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/ubuntu/express_x016/.env

# Команда запуску
ExecStart=/home/ubuntu/express_x016/.venv/bin/python -m app.main

# Автоперезапуск при падінні
# on-failure: перезапуск тільки при помилках, не при systemctl stop
Restart=on-failure
RestartSec=10

# Логування: вивід journald (див. journalctl -u express-bot -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=express-bot

[Install]
WantedBy=multi-user.target
